# tweaks/brightness-lock/zenbook-duo-brightness-lock.service
[Unit]
Description=Force Zenbook Duo lower display brightness to full
After=graphical.target
Wants=graphical.target

[Service]
Type=oneshot
Environment=TARGET_OUTPUT=eDP-2
Environment=TARGET_BRIGHTNESS=100
ExecStart=/usr/local/bin/zenbook-duo-brightness-lock.sh
TimeoutStartSec=15
